# minecraft-server/backend/panel.service
[Unit]
Description=Minecraft Control Panel API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/panel
Environment="PATH=/srv/panel/venv/bin"
ExecStart=/srv/panel/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8081
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
